# crontab
* * * * * python3 /app/redis_manager.py
@reboot python3 /app/redis_manager.py

